# initd.gentoo
#!/sbin/runscript

command=
name=$(basename $command)
pidfile=

extra_started_commands="reload reopen"

reload() {
  ebegin "Reloading ${name}"
  start-stop-daemon --signal HUP --pidfile $pidfile
  eend $?
}

reopen() {
  ebegin "Reopening ${name} log files"
  start-stop-daemon --signal USR1 --pidfile $pidfile
  eend $?
}
